# scripts/lexos-combined.service
[Unit]
Description=LexOS Combined - Frontend and Backend
After=network.target
Wants=network.target

[Service]
Type=forking
User=user
WorkingDirectory=/home/user/lexos-combined
ExecStart=/home/user/lexos-combined/scripts/auto-deploy.sh deploy
ExecStop=/home/user/lexos-combined/scripts/auto-deploy.sh stop
ExecReload=/home/user/lexos-combined/scripts/auto-deploy.sh restart
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=NODE_ENV=production
Environment=FRONTEND_PORT=3001
Environment=BACKEND_PORT=3000
Environment=HOST=0.0.0.0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/user/lexos-combined

[Install]
WantedBy=multi-user.target